# deploy/echo.service.j2
[Unit]
Description=Go echo service

[Service]
User=www-data
Group=www-data

ExecStart = /tmp/echo-service --listen=:%i

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target